# launchd/com.vadim.coinbase-crypto-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.vadim.coinbase-crypto-bot</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/vadimkoenen/Documents/Claude/Projects/Investing/alpaca-autonomous-microbot/.venv/bin/python3</string>
      <string>/Users/vadimkoenen/Documents/Claude/Projects/Investing/alpaca-autonomous-microbot/main.py</string>
      <string>--mode</string>
      <string>live</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vadimkoenen/Documents/Claude/Projects/Investing/alpaca-autonomous-microbot</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>BROKER</key>
      <string>coinbase</string>
      <key>CONFIG_FILE</key>
      <string>config_coinbase_crypto.yaml</string>
      <key>PATH</key>
      <string>/Users/vadimkoenen/Documents/Claude/Projects/Investing/alpaca-autonomous-microbot/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/vadimkoenen/Documents/Claude/Projects/Investing/alpaca-autonomous-microbot/logs/coinbase.launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/vadimkoenen/Documents/Claude/Projects/Investing/alpaca-autonomous-microbot/logs/coinbase.launchd.err.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
  </dict>
</plist>
